# rig-stats-monitor.cron
* * * * * bash ~/gpu.network/rig-stats-monitor.sh 
